# openclaw-extensions.service
[Unit]
Description=OpenClaw UI Extensions (Workbench + Channel Manager)
After=network.target
Wants=network.target

[Service]
Type=forking
User=claw-agentbox
WorkingDirectory=/media/claw-agentbox/data/9999_LocalRepo/Openclaw-OpenUSDGoodtstart-Extension
ExecStart=/media/claw-agentbox/data/9999_LocalRepo/Openclaw-OpenUSDGoodtstart-Extension/start-extension.sh start
ExecStop=/media/claw-agentbox/data/9999_LocalRepo/Openclaw-OpenUSDGoodtstart-Extension/start-extension.sh stop
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
